# imaging-from-packer/postinstall.sh
#!/bin/bash

echo "################Starting Postinstall####################"

# # Enable Open-SSH server
# systemctl enable sshd
# systemctl start sshd

echo "###################Finish Postinstall ###################"